# deploy/alfred-coo.service
[Unit]
Description=Headless Alfred COO daemon (alfred-coo-svc)
Documentation=https://github.com/salucallc/alfred-coo-svc
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=alfredcoo
Group=alfredcoo
WorkingDirectory=/opt/alfred-coo
EnvironmentFile=/etc/alfred-coo/.env
# Sub #62 (2026-04-27): point the daemon at the canonical model registry
# location on Oracle. Falls through to /opt/alfred-coo/model_registry/
# registry.yaml then to the minipc planning path if MODEL_REGISTRY_PATH
# is unset or points at a missing file. Edits to the YAML take effect on
# the NEXT dispatch — no daemon restart required.
Environment=MODEL_REGISTRY_PATH=/etc/alfred-coo/model_registry.yaml
ExecStart=/opt/alfred-coo/venv/bin/python -m alfred_coo
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alfred-coo

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/alfred-coo /var/lib/alfred-coo

[Install]
WantedBy=multi-user.target
